# systemd/sync-rhythmbox-playlists.path
[Unit]
Description="Monitor playlists file for changes"

[Path]
# we can't do recursive directory monitors (as of systemd 246) so list manually
PathModified=%h/Dropbox/Playlists/
PathModified=%h/Dropbox/Playlists/WIP
Unit=sync-rhythmbox-playlists.service

[Install]
WantedBy=default.target
